# agents/caddy-agent/caddy-manager.service
[Unit]
Description=Caddy Manager - Caddyfile staging and validation service
Documentation=file:///opt/caddy-manager/README.md
After=caddy.service
Wants=caddy.service

[Service]
Type=simple
User=caddy-manager
Group=caddy-manager
WorkingDirectory=/opt/caddy-manager

Environment="CADDYFILE_PATH=/etc/caddy/Caddyfile"
Environment="CADDY_MANAGER_DATA=/var/lib/caddy-manager"
Environment="CADDY_RELOAD_CMD=sudo systemctl reload caddy"
Environment="PORT=9999"
Environment="DEBUG=false"
Environment="PYTHONUNBUFFERED=1"

ExecStart=/usr/bin/python3 -u app.py

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=caddy-manager

[Install]
WantedBy=multi-user.target
